# packages/kszaq/rng-tools/system.d/rng-tools.service
[Unit]
Description=RNG-Tools service

[Service]
Type=oneshot
ExecStart=/usr/sbin/rngd
LimitNPROC=1
TimeoutStopSec=1s

[Install]
WantedBy=system.target
