# packages/network/btsync/system.d/btsync.service
[Unit]
Description=BitTorrent Sync Service
After=network.target

[Service]
User=btsync
ExecStart=/usr/bin/btsync --nodaemon --config /etc/btsync.conf

[Install]
WantedBy=multi-user.target
