# projects/Amlogic/devices/S859/filesystem/usr/lib/udev/rules.d/80-bluetooth-uart.rules
KERNEL=="ttyS2", SUBSYSTEMS=="tty", ENV{SYSTEMD_WANTS}="bluetooth.service"


KERNEL=="ttyS2", SUBSYSTEMS=="tty", RUN+="modprobe hci-uart"
